# systemd/agent.service
[Unit]
Description=Zero-Trust Agent
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
EnvironmentFile=/etc/agent/agent.conf
LoadCredential=CONTROLLER_CA:/etc/agent/ca.crt

Restart=always
RestartSec=10

ExecStart=/usr/bin/agent run

# The agent needs CAP_NET_ADMIN to manage nftables firewall rules.
# DynamicUser is not used because nftables requires real root-level
# netlink access that is incompatible with user namespacing.
PrivateTmp=yes
PrivateDevices=yes

ProtectSystem=full
ProtectHome=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectHostname=yes
ProtectClock=yes

AmbientCapabilities=CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN
DevicePolicy=closed

RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
RestrictNamespaces=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
SystemCallArchitectures=native
SystemCallFilter=~@clock @debug @module @mount @raw-io @reboot @cpu-emulation @obsolete

MemoryDenyWriteExecute=yes
LockPersonality=yes

UMask=0077
RuntimeDirectory=agent
StateDirectory=agent

[Install]
WantedBy=multi-user.target
