# systemd/connector.service
[Unit]
Description=Zero-Trust Connector
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
EnvironmentFile=/etc/connector/connector.conf
Environment=CONNECTOR_API_ENDPOINT=/run/connector/connector.sock
LoadCredential=CONTROLLER_CA:/etc/connector/ca.crt

Restart=always
RestartSec=10

ExecStart=/usr/bin/connector run --systemd-watchdog

DynamicUser=yes
PrivateUsers=yes
PrivateTmp=yes
PrivateDevices=yes
NoNewPrivileges=yes

ProtectSystem=full
ProtectHome=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
ProtectHostname=yes
ProtectClock=yes

CapabilityBoundingSet=
DevicePolicy=closed

RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
RestrictNamespaces=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
SystemCallArchitectures=native
SystemCallFilter=~@clock @debug @module @mount @raw-io @reboot @privileged @resources @cpu-emulation @obsolete

MemoryDenyWriteExecute=yes
LockPersonality=yes

UMask=0077
RuntimeDirectory=connector
StateDirectory=connector
Environment=STATE_DIRECTORY=/var/lib/connector
NotifyAccess=main
WatchdogSec=30s

[Install]
WantedBy=multi-user.target
